# live-build/config/includes.chroot/var/lib/dpkg/info/casper.list
/.
/etc
/etc/casper.conf
/sbin
/sbin/casper-getty
/sbin/casper-login
/sbin/casper-new-uuid
/sbin/casper-snapshot
/sbin/casper-stop
/usr
/usr/lib
/usr/lib/casper
/usr/lib/casper/casper-md5check
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/casper-md5check.service
/usr/lib/systemd/system/casper.service
/usr/share
/usr/share/casper
/usr/share/casper/casper-preseed
/usr/share/casper/casper-reconfigure
/usr/share/casper/casper-set-selections
/usr/share/casper/casper-update-initramfs
/usr/share/doc
/usr/share/doc/casper
/usr/share/doc/casper/changelog.gz
/usr/share/doc/casper/copyright
/usr/share/initramfs-tools
/usr/share/initramfs-tools/conf.d
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/casper
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/casper
/usr/share/initramfs-tools/scripts/casper-bottom
/usr/share/initramfs-tools/scripts/casper-bottom/05mountpoints
/usr/share/initramfs-tools/scripts/casper-bottom/12fstab
/usr/share/initramfs-tools/scripts/casper-bottom/13swap
/usr/share/initramfs-tools/scripts/casper-bottom/14locales
/usr/share/initramfs-tools/scripts/casper-bottom/15autologin
/usr/share/initramfs-tools/scripts/casper-bottom/18hostname
/usr/share/initramfs-tools/scripts/casper-bottom/19keyboard
/usr/share/initramfs-tools/scripts/casper-bottom/20xconfig
/usr/share/initramfs-tools/scripts/casper-bottom/22desktop_settings
/usr/share/initramfs-tools/scripts/casper-bottom/22sslcert
/usr/share/initramfs-tools/scripts/casper-bottom/24preseed
/usr/share/initramfs-tools/scripts/casper-bottom/25adduser
/usr/share/initramfs-tools/scripts/casper-bottom/25configure_init
/usr/share/initramfs-tools/scripts/casper-bottom/25disable_cdrom.mount
/usr/share/initramfs-tools/scripts/casper-bottom/30accessibility
/usr/share/initramfs-tools/scripts/casper-bottom/31disable_update_notifier
/usr/share/initramfs-tools/scripts/casper-bottom/33enable_apport_crashes
/usr/share/initramfs-tools/scripts/casper-bottom/34disable_kde_services
/usr/share/initramfs-tools/scripts/casper-bottom/37disable_screensaver_lubuntu
/usr/share/initramfs-tools/scripts/casper-bottom/40install_driver_updates
/usr/share/initramfs-tools/scripts/casper-bottom/41apt_cdrom
/usr/share/initramfs-tools/scripts/casper-bottom/43disable_updateinitramfs
/usr/share/initramfs-tools/scripts/casper-bottom/44pk_allow_ubuntu
/usr/share/initramfs-tools/scripts/casper-bottom/45jackd2
/usr/share/initramfs-tools/scripts/casper-bottom/52gnome_initial_setup
/usr/share/initramfs-tools/scripts/casper-bottom/53disable_unattended_upgrades
/usr/share/initramfs-tools/scripts/casper-bottom/55disable_snap_refresh
/usr/share/initramfs-tools/scripts/casper-bottom/56override_nvidia_udev_rule
/usr/share/initramfs-tools/scripts/casper-bottom/57pollinate
/usr/share/initramfs-tools/scripts/casper-bottom/58server_network
/usr/share/initramfs-tools/scripts/casper-bottom/59disable_mozc_autosetup
/usr/share/initramfs-tools/scripts/casper-bottom/60_create_intaller_logdir
/usr/share/initramfs-tools/scripts/casper-bottom/61desktop_canary_tweaks
/usr/share/initramfs-tools/scripts/casper-bottom/99casperboot
/usr/share/initramfs-tools/scripts/casper-functions
/usr/share/initramfs-tools/scripts/casper-helpers
/usr/share/initramfs-tools/scripts/casper-premount
/usr/share/initramfs-tools/scripts/casper-premount/10driver_updates
/usr/share/initramfs-tools/scripts/casper-premount/20iso_scan
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/casper-snapshot.1.gz
/usr/share/man/man7
/usr/share/man/man7/casper.7.gz
